# proxy/deploy/hannah-proxy.service
[Unit]
Description=Hannah Proxy — Satellite UDP/gRPC Bridge
Documentation=https://dev.kernstock.net/gessinger/voice/hannah
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/hannah-proxy --config /etc/hannah-proxy/config.yaml
Restart=on-failure
RestartSec=5s

# Runs as the shared "hannah" service user (same as all other Hannah services
# except autodeploy, which needs root to manage other systemd units).
User=hannah

# Harden the service — the proxy only needs network + its config file.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadOnlyPaths=/etc/hannah-proxy

[Install]
WantedBy=multi-user.target
